# multi-proxy-manager/launchd/com.multi-proxy-manager.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.multi-proxy-manager</string>
    <key>ProgramArguments</key>
    <array>
        <string>bash</string>
        <string>%LAUNCHD_SCRIPT_DIR%/manage.sh</string>
        <string>start</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/multi-proxy-manager-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/multi-proxy-manager-stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
